# monit-collector.ini
[program:monit_collector_main_app]
command = /opt/monit-collector/app/django-monit-collector/compose/django/gunicorn.sh                    ; Command to start app
user = root                                                          ; User to run as
stdout_logfile = /var/log/vantrix/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding



[program:monit_collector_heartbeat]
command = /opt/monit-collector/app/django-monit-collector/compose/django/heartbeat.sh                    ; Command to start app
user = root                                                          ; User to run as
stdout_logfile = /var/log/vantrix/supervisorctl_heartbeat.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding


